# dot_config/windows/scheduled-tasks/FanControl.xml
<?xml version="1.0" encoding="UTF-16"?>

<Task version="1.4" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">

  <RegistrationInfo>

    <Date>2025-10-27T09:52:40.559-07:00</Date>

    <Author>HAL9K-ALT\Randall</Author>

    <Description>Start FanControl at startup</Description>

    <URI>\FanControl</URI>

  </RegistrationInfo>

  <Principals>

    <Principal id="Author">

      <UserId>S-1-5-21-459282282-1938389383-1025489529-1000</UserId>

      <LogonType>InteractiveToken</LogonType>

      <RunLevel>HighestAvailable</RunLevel>

    </Principal>

  </Principals>

  <Settings>

    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>

    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>

    <ExecutionTimeLimit>PT0S</ExecutionTimeLimit>

    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>

    <Priority>3</Priority>

    <IdleSettings>

      <StopOnIdleEnd>true</StopOnIdleEnd>

      <RestartOnIdle>false</RestartOnIdle>

    </IdleSettings>

    <UseUnifiedSchedulingEngine>true</UseUnifiedSchedulingEngine>

  </Settings>

  <Triggers>

    <LogonTrigger>

      <UserId>HAL9K-ALT\Randall</UserId>

    </LogonTrigger>

  </Triggers>

  <Actions Context="Author">

    <Exec>

      <Command>FanControl</Command>

      <WorkingDirectory>C:\Users\Randall\Downloads\FanControl_245_net_8_0\</WorkingDirectory>

    </Exec>

  </Actions>

</Task>
